# deployment/pastebin-app.service
[Unit]
Description=Pastebin WSGI Application
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/pastebin
Environment="PATH=/var/www/pastebin/venv/bin"

# Use gunicorn for production (better than wsgiref)
ExecStart=/var/www/pastebin/venv/bin/gunicorn \
    --workers 4 \
    --bind 0.0.0.0:8000 \
    --access-logfile /var/www/pastebin/logs/access.log \
    --error-logfile /var/www/pastebin/logs/error.log \
    --log-level info \
    backend.app:application

# Restart policy
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/pastebin/logs

# Performance tuning
LimitNOFILE=65535

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pastebin-app

[Install]
WantedBy=multi-user.target